# scripts/mt.service
# ============================================================================
# MT - Systemd Service File
# ============================================================================
# Resurrects PM2 managed processes after system reboot.
#
# Install:
#   sudo cp /root/scripts/mt.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable mt
#   sudo systemctl start mt
#
# Verify:
#   sudo systemctl status mt
#   pm2 list
#
# Logs:
#   journalctl -u mt -f
# ============================================================================

[Unit]
Description=MT - PM2 Process Manager
Documentation=https://github.com/your-org/mt
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=forking
User=root
Group=root
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/root/.local/share/pnpm:/root/.nvm/versions/node/v18.*/bin
Environment=PM2_HOME=/root/.pm2
Environment=PROJECT_ROOT=/root
PIDFile=/root/.pm2/pm2.pid

# Start: resurrect PM2 saved process list
ExecStart=/usr/local/bin/pm2 resurrect
# Reload: graceful reload of all processes
ExecReload=/usr/local/bin/pm2 reload all
# Stop: graceful shutdown of all processes
ExecStop=/usr/local/bin/pm2 kill

# Restart on failure
Restart=on-failure
RestartSec=10

# Hardening
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mt

[Install]
WantedBy=multi-user.target
